# config/statsd/statsd.conf
[program:statsd]
command=/usr/local/bin/node /opt/statsd/stats.js /opt/statsd/localConfig.js
autostart=true
autorestart=true
stopsignal=QUIT
user=vagrant
